# modules/b2g_bumper/templates/run_b2g_bumper.cron.erb
# This file is managed by puppet
*/5 * * * *	<%=@b2g_bumper_user%>	/usr/local/bin/run_b2g_bumper.sh
