# modules/puppet/files/puppet.upstart.conf
description "start puppet"

# Don't expect long running process
task

start on (local-filesystems and net-device-up IFACE!=lo)
exec /etc/puppet/init start
