# modules/nightlypromotion/templates/run_nightlypromotion.cron.erb
# This file is managed by puppet
0 10 * * *	<%=@nightlypromotion_user%>		<%=scope.lookupvar('nightlypromotion::settings::root')%>/bin/python2.7 <%=scope.lookupvar('nightlypromotion::settings::script')%> <%=scope.lookupvar('nightlypromotion::settings::authfile')%> <%=scope.lookupvar('nightlypromotion::settings::cacert')%> --log-file <%=scope.lookupvar('nightlypromotion::settings::logfile')%> --cache-file <%=scope.lookupvar('nightlypromotion::settings::cachefile')%>
