# modules/needs_reboot/files/rm_reboot.upstart.conf
description "remove /REBOOT_AFTER_PUPPET"

task

start on (local-filesystems)

exec /bin/rm -rf /REBOOT_AFTER_PUPPET
