# modules/buildslave/templates/linux-initd-buildbot.sh.erb
#! /bin/bash
#  initscript for buildbot
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

# chkconfig: 2345 99 08
### BEGIN INIT INFO
# Provides:          buildbot
# Required-Start:    $local_fs $network puppet runner supervisord
# Should-Start:      $remote_fs 
# Should-Stop:       $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      K 0 1 6
# Short-Description: Buildbot
# Description:       Buildbot
### END INIT INFO

# MOZILLA DEPLOYMENT NOTES
# - This file is distributed by Puppet to all linux build servers
# - it lives in the 'buildslave' puppet module

# Required-Stop:     $local_fs
PATH=/sbin:/bin:/usr/sbin:/usr/bin
DESC="Buildbot"

# this uses the /tools/buildbot symlink to use the active version of buildbot
BUILDSLAVE_CMD=/tools/buildbot/bin/buildslave
LOCKFILE=/var/lock/subsys/buildbot
RUNSLAVE=/usr/local/bin/runslave.py
# use the python from the buildbot virtualenv
PYTHON=/tools/buildbot/bin/python

test -x ${BUILDSLAVE_CMD} || exit 0

. /etc/init.d/functions

start_buildbot() {
    # note - spaces here will break things
    su - <%= scope.lookupvar('::config::builder_username') %> sh -c "${PYTHON} ${RUNSLAVE}"
    ret=$?
    if [ $ret == "0" ]; then
        touch ${LOCKFILE}
    fi
    return $ret
}

do_start () {
    errors=0
	echo "Starting buildslave"
	if start_buildbot 
	then
	    echo "started"
	else
	    echo "not started"
	    errors=$(($errors+1))
	fi
    return $errors
}

do_stop () {
    echo "sorry, no-can-do.  This script is a one-way-trip to start the buildslave."
    return 1
}

do_reload () {
    do_stop
}

do_restart () {
    do_stop
}

case "$1" in
  start)
  	do_start
  	exit $?
	;;
  stop)
  	do_stop
  	exit $?
	;;
  reload)
  	do_reload
  	exit $?
	;;
  restart|force-reload)
  	do_restart
  	exit $?
	;;
  *)
	log_warning_msg "Usage: $0 {start|stop|restart|reload|force-reload}"
	exit 1
	;;
esac

exit 0
