# modules/buildmaster/templates/pulse_publisher.initd.erb
#!/bin/bash
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.
# Start up the pulse_publisher daemon

# chkconfig: 2345 50 08
### BEGIN INIT INFO
# Provides:          pulse_publisher
# Required-Start:    $local_fs $network
# Should-Start:      $remote_fs
# Should-Stop:       $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      K 0 1 6
# Short-Description: pulse publisher
# Description:       pulse publisher
### END INIT INFO

# Required-Stop:     $local_fs
PATH=/sbin:/bin:/usr/sbin:/usr/bin
DESC="pulse publisher"

USER=<%=scope.lookupvar('users::builder::username')%>
BASE_DIR=<%=scope.lookupvar('buildmaster::settings::queue_dir')%>
PYTHON="$BASE_DIR/bin/python"
SCRIPT="$BASE_DIR/tools/buildbot-helpers/pulse_publisher.py"
LOGFILE="$BASE_DIR/pulse_publisher.log"
LOCKFILE="/var/lock/subsys/pulse_publisher"
PIDFILE="$BASE_DIR/pulse_publisher.pid"
OPTIONS="--passwords $BASE_DIR/passwords.py -q /dev/shm/queue/pulse >> $LOGFILE"


test -x "${SCRIPT}" || exit 0

. /etc/init.d/functions

start() {
    echo -n "Starting $DESC"
    daemon --user $USER --pidfile $PIDFILE "$PYTHON $SCRIPT $OPTIONS &"
    RETVAL=$?
    PID=$(ps -f -U "$USER" | grep -F "$SCRIPT" | awk {'print $2'})
    if [ -n "$PID" ]; then
        echo $PID > "$PIDFILE"
        touch "$LOCKFILE"
    fi
    echo
    return $RETVAL
}

stop() {
    echo -n "Stopping $DESC"
    killproc -p "$PIDFILE" -d 1 "$SCRIPT"
    RETVAL=$?
    [ $RETVAL -eq 0 ] && rm -rf "$LOCKFILE"
    echo
    return $RETVAL
}

restart () {
    stop
    start
    RETVAL=$?
    echo
    return $RETVAL
}

do_status () {
    echo -n "$DESC"
    status -p $PIDFILE
    return $?
}

case "$1" in
  start)
        start
        exit $?
        ;;
  stop)
        stop
        exit $?
        ;;
  restart)
        restart
        exit $?
        ;;
  status)
        do_status
        exit $?
        ;;
  *)
        log_warning_msg "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0
