# modules/puppetmaster/templates/foreman_facts.cron.erb
MAILTO="<%= scope.lookupvar('::config::puppet_notif_email') %>"
*/30 * * * *     root    (sleep $(($RANDOM/(32768/30)))); <%= @script %> 2>&1 | logger -t foreman_fact_import
MAILTO=root
