# modules/puppetmaster/templates/ssl_git_sync.cron.erb
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.
MAILTO="<%= scope.lookupvar('::config::puppet_notif_email') %>"
*/30 * * * * puppetsync (sleep $(($RANDOM\%1800))) && <%= @scripts_dir %>/ssl_git_sync.sh
MAILTO=root
