# modules/gaia_bumper/templates/run_gaia_bumper.cron.erb
# This file is managed by puppet
*/5 * * * *	<%=@gaia_bumper_user%>	/usr/local/bin/run_gaia_bumper.sh
