# modules/puppet/templates/puppet-centos-initd.erb
#!/bin/bash
# puppet        Init script for running the puppet client daemon
#
# Author:       Duane Griffin <d.griffin@psenterprise.com>
#               David Lutterkort <dlutter@redhat.com>
#
# chkconfig: 2345 97 02
### BEGIN INIT INFO
# Provides:          puppet
# Required-Start:    $local_fs $network
# Required-Start:    $local_fs $network
# Should-Start:      $remote_fs
# Should-Stop:       $remote_fs
# Default-Start:     2 3 4 5
#
# description: Enables periodic system configuration checks through puppet.
# processname: puppet
# config: /etc/sysconfig/puppet

PATH=/usr/bin:/sbin:/bin:/usr/sbin
export PATH

# Source function library.
. /etc/rc.d/init.d/functions

# common atboot code
<%= @puppet_atboot_common %>

case "$1" in
  start)
    # -v to show that we're removing the lock in logs
    rm -vf /var/lib/puppet/state/agent_catalog_run.lock
    run_until_success
    ;;
  stop)
    exit 0 # sure, I stopped..
    ;;
  status)
    exit 0 # working just fine, thanks
    ;;
  *)
    echo $"Usage: $0 {start|stop|status}"
    exit 1
esac
