# modules/runner/templates/runner.service.erb
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.
#
#  initscript for runner

[Unit]
Description=Job that runs the runner daemon
Documentation=

[Service]
PIDFile=/var/run/runner.pid
ExecStart=/usr/bin/nohup /opt/runner/bin/runner -v --syslog -n 30 -H -c /opt/runner/runner.conf /opt/runner/tasks.d

[Install]
WantedBy=multi-user.target