# modules/runner/templates/runner.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<%# This Source Code Form is subject to the terms of the Mozilla Public
  # License, v. 2.0. If a copy of the MPL was not distributed with this
  # file, You can obtain one at http://mozilla.org/MPL/2.0/. -%>

 <!-- note that this service is run by the user-level launchd, not by the
      system launchd on login, so it does not need to specify a UserName (and in
      fact, this field is ignored) -->

<dict>
        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin:/usr/X11/bin</string>
        </dict>
        <key>Label</key>
        <string>com.mozilla.runner</string>
        <key>ProgramArguments</key>
        <array>
            <string><%= scope.lookupvar('::runner::settings::root') %>/bin/runner</string>
            <string>-v</string>
            <string>--syslog</string>
            <string>-H</string>
            <string>-n</string>
            <string>15</string>
            <string>-c</string>
            <string><%= scope.lookupvar('::runner::settings::root') %>/runner.cfg</string>
            <string><%= scope.lookupvar('::runner::settings::root') %>/tasks.d</string>
        </array>

        <!-- do not run immediately when loaded -->
        <key>RunAtLoad</key>
        <false/>

        <!-- but run when this file exists -->
        <key>KeepAlive</key>
        <dict>
            <key>PathState</key>
            <dict>
                <!-- start if this file exists -->
                <key>/var/tmp/semaphore/run-buildbot</key>
                <true/>
            </dict>
        </dict>

        <key>OnDemand</key>
        <true/>

        <!-- the working directory doesn't get used, so point it somewhere inoccuous -->
        <key>WorkingDirectory</key>
        <string>/Users</string>

        <key>StandardErrorPath</key>
        <string>/var/tmp/runner.err</string>

        <key>StandardOutPath</key>
        <string>/var/tmp/runner.out</string>
</dict>
</plist>
