# modules/runner/templates/runner.initd.erb
#! /bin/bash
#  initscript for runner
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

# chkconfig: 2345 45 08
### BEGIN INIT INFO
# Provides:          runner
# Required-Start:    $local_fs $network puppet
# Should-Start:      $remote_fs
# Should-Stop:       $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      K 0 1 6
# Short-Description: runner
# Description:       runs tasks
### END INIT INFO

# MOZILLA DEPLOYMENT NOTES
# - This file is distributed by Puppet to all linux build servers
# - it lives in the 'runner' puppet module

# Required-Stop:     $local_fs
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC="runner"

CMD=<%= scope.lookupvar("::runner::settings::root") %>/bin/runner
CONFIG=<%= scope.lookupvar("::runner::settings::root") %>/runner.cfg
TASKDIR=<%= scope.lookupvar("::runner::settings::taskdir") %>

test -x ${CMD} || exit 0

do_start () {
    echo "Starting runner"
    ### XXX -n 1 for testing
    if ${CMD} -c ${CONFIG} -n 1 ${TASKDIR}; then
        echo "done"
        return 0
    else
        echo "failed"
        return 1
    fi
}

do_stop () {
    echo "nothing to do"
    return 0
}

do_reload () {
    do_stop
}

do_restart () {
    do_stop
}

case "$1" in
  start)
        do_start
        exit $?
        ;;
  stop)
        do_stop
        exit $?
        ;;
  reload)
        do_reload
        exit $?
        ;;
  restart|force-reload)
        do_restart
        exit $?
        ;;
  *)
        echo "Usage: $0 {start|stop|restart|reload|force-reload}"
        exit 1
        ;;
esac

exit 0
