# modules/runner/templates/runner.upstart.conf.erb
description "start runner"

start on stopped puppet and (started x11 or started xvfb)

script
    CMD=<%= scope.lookupvar("::runner::settings::root") %>/bin/runner
    CONFIG=<%= scope.lookupvar("::runner::settings::root") %>/runner.cfg
    TASKDIR=<%= scope.lookupvar("::runner::settings::taskdir") %>

    # We sleep a bit here because even though Xvfb has completed, we want to
    # make sure that the DE has launched. Some sort of check of the process
    # list would be better, but this is probably good enough.
    sleep 30

    /usr/sbin/logrotate -f /etc/logrotate.d/runner
    ${CMD} -v --syslog -n 30 -H -c ${CONFIG} ${TASKDIR} >> /var/log/runner.log 2>&1
end script
