# modules/vnc/templates/x11vnc.service.erb
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.
#
# Define x11vnc service

[Unit]
Description=Manually start x11vnc server attached to <%= scope.lookupvar('::config::builder_username') %> X session
After=Xsession.service

[Service]
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
# this file is stored in /etc, readable only by root.  When starting VNC, it gets copied into
# the builder user's home dir, then removed immediately (via the 'rm:' prefix below)
ExecStartPre=/bin/cp -f /etc/vnc_passwdfile /home/<%= scope.lookupvar('::config::builder_username') %>/.vnc/passwdfile
ExecStartPre=/bin/chown <%= scope.lookupvar('::config::builder_username') %> /home/<%= scope.lookupvar('::config::builder_username') %>/.vnc/passwdfile
ExecStart=/bin/su - -c "/usr/bin/x11vnc -display WAIT:0 -passwdfile rm:/home/<%= scope.lookupvar('::config::builder_username') %>/.vnc/passwdfile -forever -shared" <%= scope.lookupvar('::config::builder_username') %>

[Install]
WantedBy=multi-user.target