# modules/vnc/templates/xvfb.conf.erb
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

description "start Xvfb server"

export USER=<%= scope.lookupvar('::config::builder_username') %>
setuid <%= scope.lookupvar('::config::builder_username') %>


# wait until puppet finishes
start on stopped puppet

exec Xvfb :0 -screen 0 <%= scope.lookupvar('::screenresolution::talos::width') %>x<%= scope.lookupvar('::screenresolution::talos::height') %>x<%= scope.lookupvar('::screenresolution::talos::xvfb_depth') %>
