# modules/instance_metadata/templates/instance_metadata.initd.erb
#! /bin/bash
#  initscript for instance_metadata
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

# chkconfig: 2345 49 08
### BEGIN INIT INFO
# Provides:          instance_metadata
# Required-Start:    $local_fs $network buildbot
# Should-Start:      $remote_fs
# Should-Stop:       $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      K 0 1 6
# Short-Description: instance_metadata
# Description:       Fetches instance metadata for VMs and stores it in a local file
### END INIT INFO

# MOZILLA DEPLOYMENT NOTES
# - This file is distributed by Puppet to all linux build servers
# - it lives in the 'instance_metadata' puppet module

# Required-Stop:     $local_fs
PATH=/sbin:/bin:/usr/sbin:/usr/bin
DESC="instance_metadata"

CMD=/usr/local/bin/instance_metadata.py
OUTPUT=/etc/instance_metadata.json
PYTHON=<%= scope.lookupvar('::packages::mozilla::python27::python') %>

test -x ${CMD} || exit 0

do_start () {
    echo "Fetching instance metadata"
    if ${CMD} -o ${OUTPUT}; then
        echo "done"
        return 0
    else
        echo "failed"
        return 1
    fi
}

do_stop () {
    echo "nothing to do"
    return 1
}

do_reload () {
    do_stop
}

do_restart () {
    do_stop
}

case "$1" in
  start)
  	do_start
  	exit $?
	;;
  stop)
  	do_stop
  	exit $?
	;;
  reload)
  	do_reload
  	exit $?
	;;
  restart|force-reload)
  	do_restart
  	exit $?
	;;
  *)
	log_warning_msg "Usage: $0 {start|stop|restart|reload|force-reload}"
	exit 1
	;;
esac

exit 0
