# modules/gui/templates/Xsession.conf.erb
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

description "Start <%= scope.lookupvar('::config::builder_username') %> X session in <%= if @on_gpu then 'x11' else 'xvfb' end %>"

start on started <%= if @on_gpu then 'x11' else 'xvfb' end %>

script
  # note that doing anything with this display before running Xsession will cause sadness (like,
  # crashes in compiz)
  # Make sure that X has enough time to start
  sleep 15
  su - -c "DISPLAY=:0 /etc/X11/Xsession" <%= scope.lookupvar('::config::builder_username') %>
end script
