# modules/gui/templates/xvfb.conf.erb
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

description "start Xvfb server"

export USER=<%= scope.lookupvar('::config::builder_username') %>
setuid <%= scope.lookupvar('::config::builder_username') %>


# wait until puppet finishes
start on stopped puppet

<%# note that xvfb's maximum depth is 24 -%>
exec Xvfb :0 -nolisten tcp -screen 0 <%= @screen_width %>x<%= @screen_height %>x<%= if @screen_depth.to_i <= 24 then @screen_depth else 24 end %>
